# contrib/init/INTICOINd.conf
description "INTICOIN Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env INTICOIND_BIN="/usr/bin/INTICOINd"
env INTICOIND_USER="INTICOIN"
env INTICOIND_GROUP="INTICOIN"
env INTICOIND_PIDDIR="/var/run/INTICOINd"
# upstart can't handle variables constructed with other variables
env INTICOIND_PIDFILE="/var/run/INTICOINd/INTICOINd.pid"
env INTICOIND_CONFIGFILE="/etc/INTICOIN/INTICOIN.conf"
env INTICOIND_DATADIR="/var/lib/INTICOINd"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # INTICOINd will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$INTICOIND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run INTICOINd."
        echo "The setting must appear in $INTICOIND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"INTICOIN Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$INTICOIND_PIDDIR"
    chmod 0755 "$INTICOIND_PIDDIR"
    chown $INTICOIND_USER:$INTICOIND_GROUP "$INTICOIND_PIDDIR"
    chown $INTICOIND_USER:$INTICOIND_GROUP "$INTICOIND_CONFIGFILE"
    chmod 0660 "$INTICOIND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$INTICOIND_PIDFILE" \
    --chuid $INTICOIND_USER:$INTICOIND_GROUP \
    --exec "$INTICOIND_BIN" \
    -- \
    -pid="$INTICOIND_PIDFILE" \
    -conf="$INTICOIND_CONFIGFILE" \
    -datadir="$INTICOIND_DATADIR" \
    -disablewallet \
    -daemon

